# deploy/incubator.service
[Unit]
Description=Incubator v3 (Pi Zero 2W)
# Start the app even if the network isn't fully up yet — it manages its own
# Wi-Fi hotspot for onboarding, so it must not wait on network-online at first
# boot (there is no network then). NetworkManager just needs to be present.
After=NetworkManager.service incubator-firstboot.service
Wants=NetworkManager.service
# Appliance behaviour: never give up restarting. Without this, systemd's default
# start-rate limit (5 starts / 10 s) would stop the service after a brief boot
# crash-loop and leave the incubator dead until a manual power cycle.
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=/etc/incubator.env
ExecStart=__INSTALL_DIR__/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --log-level info
Restart=always
RestartSec=5
# Allow GPIO and serial access
SupplementaryGroups=gpio dialout video

[Install]
WantedBy=multi-user.target
