# deploy/incubator-v3.service
[Unit]
Description=Incubator v3 API (UNO Q)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=/etc/incubator-v3.env
ExecStart=__INSTALL_DIR__/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
